# play2/templates/default/app_initd.erb
#!/bin/sh
# chkconfig: 345 20 80
# description: <%= @name %> playframework application start/shutdown script
# processname: <%= @name %>

# User running the Play process
USER=root

# Name of the application
APPLICATION_NAME="<%= @name %>"
# Path to the application source
APPLICATION_PATH="<%= @path %>"
# Command line options for starting Play
APPLICATION_OPTIONS="<%= @options %>"
# URL to ping the active server
APPLICATION_URL="<%= @url %>"

. /etc/init.d/functions

RETVAL=0

[ -d $APPLICATION_PATH/logs ] || exit 1
[ -f $APPLICATION_PATH/<%= @command %> ] || exit 1

[ -f $APPLICATION_PATH/logs/service.log ] || touch $APPLICATION_PATH/logs/service.log

start() {
    echo -n "Starting Play service: <%= @name %>" | tee -a $APPLICATION_PATH/logs/service.log
    echo >> $APPLICATION_PATH/logs/service.log

    CHECK_UP=<<CMD
        RETVAL=$?
        PID=$!
        if [ $RETVAL -eq 0 ]; then
            while kill -0 $PID 2> /dev/null; do
                echo -n "Pinging $APPLICATION_URL" | tee -a $APPLICATION_PATH/logs/service.log
                curl -k $APPLICATION_URL > /dev/null 2> /dev/null
                if [ $? -eq 0 ]; then
                    exit 0
                fi
                sleep 1.5
            done
            exit 1
        else
            exit $RETVAL
        fi
CMD
    su $USER -c "exec ${APPLICATION_PATH}/<%= @command %> ${APPLICATION_OPTIONS} >> $APPLICATION_PATH/logs/service.log 2>&1 & $CHECK_UP"

    RETVAL=$?

    if [ $RETVAL -eq 0 ]; then
        echo_success
    else
        echo_failure
    fi
    echo

    return $RETVAL
}

stop() {
    echo -n "Shutting down Play service: <%= @name %>" | tee -a $APPLICATION_PATH/logs/service.log
    echo >> $APPLICATION_PATH/logs/service.log

    for PID in `find <%= @deploy_to %>/releases -name RUNNING_PID`
    do
        echo "Found $PID" >> $APPLICATION_PATH/logs/service.log
        if ps -p `cat $PID` > /dev/null
        then
           su $USER -c "kill -9 `cat $PID`"
        fi

        su $USER -c "rm -rf $PID"
    done

    echo_success
    echo
}

status() {
    #TODO
    echo
}

clean() {
    echo "Cleaning up Play service: <%= @name %>" >> $APPLICATION_PATH/logs/service.log
    rm -f ${APPLICATION_PATH}/RUNNING_PID
}

case "$1" in
    start)
    clean
    start
    ;;
    stop)
    stop
    ;;
    restart|reload)
    stop
    sleep 10
    start
    ;;
    status)
    status
    ;;
    clean)
    clean
    ;;
    *)
    echo "Usage: $0 {start|stop|restart|status}"
esac
exit 0
