# play2/templates/default/app_initd.erb
#!/bin/sh
# chkconfig: 345 20 80
# description: <%= @name %> playframework application start/shutdown script
# processname: <%= @name %>

# User running the Play process
USER=root

# Name of the application
APPLICATION_NAME="<%= @name %>"
# Path to the application source
APPLICATION_PATH="<%= @path %>"
APP_OPTIONS="<%= @options %>"

. /etc/init.d/functions

RETVAL=0

[ -d $APPLICATION_PATH/logs ] || exit 1
[ -f $APPLICATION_PATH/<%= @command %> ] || exit 1

[ -f $APPLICATION_PATH/logs/service.log ] || touch $APPLICATION_PATH/logs/service.log

start() {
    echo -n "Starting Play service: <%= @name %>" | tee -a $APPLICATION_PATH/logs/service.log
    su $USER -c "exec ${APPLICATION_PATH}/<%= @command %> ${APP_OPTIONS} >> $APPLICATION_PATH/logs/service.log 2>&1 &"
    RETVAL=$?

    if [ $RETVAL -eq 0 ]; then
        echo_success
    else
        echo_failure
    fi
    echo
}

stop() {
    echo -n "Shutting down Play service: <%= @name %>" | tee -a $APPLICATION_PATH/logs/service.log
    for PID in `find <%= @deploy_to %>/releases -name RUNNING_PID`
    do
        echo "Found $PID" >> $APPLICATION_PATH/logs/service.log
        su $USER -c "kill -9 `cat $PID` ; rm -rf $PID"
    done

    echo_success
    echo
}

status() {
    #TODO
    echo
}

clean() {
    echo "Cleaning up Play service: <%= @name %>" >> $APPLICATION_PATH/logs/service.log
    rm -f ${APPLICATION_PATH}/RUNNING_PID
}

case "$1" in
    start)
    clean
    start
    ;;
    stop)
    stop
    ;;
    restart|reload)
    stop
    sleep 10
    start
    ;;
    status)
    status
    ;;
    clean)
    clean
    ;;
    *)
    echo "Usage: $0 {start|stop|restart|status}"
esac
exit 0
